# i2pd/debian/i2pd.upstart
description "i2p client daemon"

start on runlevel [2345]
stop on runlevel [016] or unmounting-filesystem

# these can be overridden in /etc/init/i2pd.override
env LOGFILE="/var/log/i2pd/i2pd.log"

expect fork

exec /usr/sbin/i2pd --daemon --service --log=file --logfile=$LOGFILE
